# deploy/systemd/nano-app.service
[Unit]
Description=Nano app hub (web UI/gateway + OLED + behaviour statechart, one process)
After=nano-router.service
Wants=nano-router.service
PartOf=nano-robot.target

[Service]
# Type=notify: the hub sends READY=1 once up, then pets WATCHDOG=1 every 5 s from an
# executor timer — so an alive-but-wedged executor (a stuck callback) stops petting
# and systemd restarts the hub. Crashes are Restart=on-failure as before.
Type=notify
WatchdogSec=90
User=ibster
# GPU vision (gpu_vision.py) needs the `lima` kernel module. /etc/modules-load.d/lima.conf
# (see deploy/sbc-setup.sh) is supposed to load it at boot but has proven unreliable —
# systemd-modules-load.service runs very early, before whatever the Mali GPU platform
# device depends on is ready, so `lima` can silently fail to bind even though the conf
# file is correct (confirmed on hardware 2026-07-11 — see memory lima-boot-load-bug).
# Retry it here, right before app_hub actually needs it: by the time nano-app starts
# (After=nano-router.service, itself gated on the router actually listening) far more
# of the system has settled than at early-boot modules-load time. `-` = non-fatal if it
# fails; gpu_vision.py already degrades gracefully (falls back, now loudly logs a
# warning) rather than crashing if the GPU truly isn't available. `+` = run this ONE
# command as root regardless of the unit's User=ibster below -- loading a kernel module
# needs CAP_SYS_MODULE, which ibster doesn't have; confirmed on hardware 2026-07-11
# that omitting `+` fails with "modprobe: ERROR: could not insert 'lima': Operation not
# permitted" (silently, since `-` swallows the failure) -- don't drop this prefix.
ExecStartPre=-+/usr/sbin/modprobe lima
ExecStart=/home/ibster/Nano/scripts/unit_exec.sh app
# SIGTERM lets the OLED draw its restart/shutdown end-screen before exiting.
TimeoutStopSec=15
Restart=on-failure
RestartSec=3
# A leak (LLM/cognition is the likely candidate) kills and restarts THIS hub instead
# of inviting the kernel OOM killer to pick a victim on the 1 GB board.
MemoryMax=450M

[Install]
WantedBy=nano-robot.target
