# deploy/systemd/nano-heal.service
[Unit]
Description=Nano stack autoheal (relaunch any crashed node)
# Only heal while the stack is supposed to be up: Requisite fails this unit
# (without starting it) if nano-stack.service isn't active, so a manual
# `systemctl stop nano-stack` also stops healing. The timer is gated the same way.
After=nano-stack.service
Requisite=nano-stack.service

[Service]
# stack.sh heal == idempotent `up`: relaunches only dead nodes (pgrep-guarded),
# a no-op when everything is alive. Oneshot — fires, heals, exits.
Type=oneshot
User=ibster
WorkingDirectory=/home/ibster/Nano
ExecStart=/home/ibster/.pixi/bin/pixi run bash scripts/stack.sh heal
TimeoutStartSec=120
# Same as nano-stack.service: the nodes heal relaunches are setsid-detached but stay
# in THIS unit's cgroup, and a oneshot's cgroup is reaped the moment it deactivates —
# without this, every node heal ever "revived" was silently SIGKILLed seconds later
# and relaunched on the next tick, forever (a 20 s crash-loop treadmill).
KillMode=process
