# deploy/systemd/nano-heal.timer
[Unit]
Description=Run Nano stack autoheal periodically
# Don't tick if the stack unit isn't loaded/active.
After=nano-stack.service
Requisite=nano-stack.service

[Timer]
# First check a minute after boot (let nano-stack.service finish its own startup,
# TimeoutStartSec=180), then every 20s. AccuracySec keeps the wakeups cheap on the H5.
OnBootSec=90
OnUnitActiveSec=20
AccuracySec=2s

[Install]
WantedBy=timers.target
